# cake-autorate
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org

export PATH=/usr/sbin:/usr/bin:/sbin:/bin
export HOME=/root

exec &> /tmp/cake-autorate.log

START=97
STOP=4

start() 
{
        bash /root/CAKE-autorate/CAKE-autorate.sh&
	echo "PID="$! > /tmp/CAKE-autorate-PID
}

stop()
{
	if [[ -f "/tmp/CAKE-autorate-PID" ]]; then
		. /tmp/CAKE-autorate-PID
		kill $PID
		rm "/tmp/CAKE-autorate-PID"
	fi
}
